# deploy/tracking-app.service
[Unit]
Description=Tracking Logistico - Gunicorn WSGI Server
After=network.target mariadb.service redis-server.service
Wants=mariadb.service redis-server.service

[Service]
# Type=simple: gunicorn no emite sd_notify por defecto.
# Si se instala systemd-python y se habilita el hook, cambiar a notify.
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tracking-app
EnvironmentFile=/etc/tracking-app/env

ExecStart=/opt/tracking-app/.venv/bin/gunicorn -c gunicorn_config.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID

Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

# ── Hardening systemd ──
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tracking-app/instance /var/data/tracking/evidencias /var/log/tracking
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
